# debian/postinst
#!/bin/bash
set -e

# Marker comment to identify our configuration
MULTIPATH_MARKER_BEGIN="# BEGIN jt-pve-storage-netapp"
MULTIPATH_MARKER_END="# END jt-pve-storage-netapp"
MULTIPATH_CONF="/etc/multipath.conf"

# Function to add NetApp multipath configuration
add_netapp_multipath_config() {
    # Ensure /etc exists (should always exist, but be safe)
    if [ ! -d "/etc" ]; then
        echo "ERROR: /etc directory does not exist!"
        return 1
    fi

    # Check if our configuration already exists
    if [ -f "$MULTIPATH_CONF" ] && grep -q "$MULTIPATH_MARKER_BEGIN" "$MULTIPATH_CONF" 2>/dev/null; then
        echo "NetApp multipath configuration already exists in $MULTIPATH_CONF"
        return 0
    fi

    # Check if any NetApp device configuration already exists (manual setup)
    if [ -f "$MULTIPATH_CONF" ] && grep -q 'vendor.*"NETAPP"' "$MULTIPATH_CONF" 2>/dev/null; then
        echo "Existing NetApp configuration found in $MULTIPATH_CONF, skipping auto-configuration."
        echo "Note: Manual configuration will be preserved."
        echo ""
        # Check for dangerous infinite-queue settings
        local needs_update=0
        if grep -q 'no_path_retry.*queue' "$MULTIPATH_CONF" 2>/dev/null; then
            needs_update=1
        fi
        if grep -q 'queue_if_no_path' "$MULTIPATH_CONF" 2>/dev/null; then
            needs_update=1
        fi
        if grep -q 'dev_loss_tmo.*infinity' "$MULTIPATH_CONF" 2>/dev/null; then
            needs_update=1
        fi
        if [ "$needs_update" -eq 1 ]; then
            echo ""
            echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
            echo -e "\033[1;31m  WARNING: DANGEROUS MULTIPATH SETTINGS DETECTED\033[0m"
            echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
            echo ""
            echo -e "  Your \033[1m$MULTIPATH_CONF\033[0m contains settings that may cause"
            echo -e "  the \033[1;31mENTIRE PVE NODE TO HANG\033[0m if a NetApp LUN becomes unavailable."
            echo ""
            echo -e "  \033[1mRequired changes:\033[0m"
            echo ""
            echo -e "    \033[1;31mno_path_retry queue\033[0m    -->  \033[1;32mno_path_retry 30\033[0m"
            echo -e "    \033[1;31mqueue_if_no_path\033[0m       -->  \033[1;32m(remove, use no_path_retry instead)\033[0m"
            echo -e "    \033[1;31mdev_loss_tmo infinity\033[0m  -->  \033[1;32mdev_loss_tmo 60\033[0m"
            echo ""
            echo -e "  \033[1mAlso add if not present:\033[0m"
            echo -e "    \033[1;32mfast_io_fail_tmo 5\033[0m"
            echo ""
            echo -e "  \033[1mAfter editing multipath.conf:\033[0m"
            echo -e "    \033[1;32msystemctl restart multipathd\033[0m   (applies new config and flushes stale maps)"
            echo ""
            echo -e "  \033[1;33mIMPORTANT:\033[0m 'systemctl reload multipathd' only re-reads config but does"
            echo -e "  NOT clean up existing stale multipath maps. You must use 'restart' to"
            echo -e "  apply new settings to existing devices and remove orphaned entries."
            echo ""
            echo -e "  \033[1;31mWARNING:\033[0m Do NOT use 'multipath -F' (capital F) to clean stale maps."
            echo -e "  It flushes ALL unused maps, which may include working maps from other"
            echo -e "  storage systems. Instead, target specific stale WWIDs:"
            echo -e "    \033[1;32mmultipath -f <wwid>\033[0m   (lowercase f, removes one specific map)"
            echo ""
            echo -e "  For details, see: \033[1;36mhttps://github.com/jasoncheng7115/jt-pve-storage-netapp\033[0m"
            echo -e "  Refer to \033[1mdocs/CONFIGURATION.md\033[0m -> Multipath Configuration section."
            echo ""
            echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
            echo ""
        fi
        return 0
    fi

    echo "Adding NetApp multipath configuration to $MULTIPATH_CONF..."

    # Create backup if file exists
    if [ -f "$MULTIPATH_CONF" ]; then
        cp "$MULTIPATH_CONF" "${MULTIPATH_CONF}.bak.$(date +%Y%m%d%H%M%S)"
        echo "Backup created: ${MULTIPATH_CONF}.bak.$(date +%Y%m%d%H%M%S)"
    fi

    # If multipath.conf doesn't exist, create with defaults
    if [ ! -f "$MULTIPATH_CONF" ]; then
        echo "Creating new $MULTIPATH_CONF..."
        cat > "$MULTIPATH_CONF" << 'DEFAULTS_EOF'
# Multipath configuration file
# Created by jt-pve-storage-netapp package

defaults {
    user_friendly_names yes
    find_multipaths yes
}

DEFAULTS_EOF
    fi

    # Add our configuration
    cat >> "$MULTIPATH_CONF" << 'EOF'

# BEGIN jt-pve-storage-netapp
# NetApp ONTAP device configuration
# Added automatically by jt-pve-storage-netapp package
# This section will be removed when the package is purged
devices {
    device {
        vendor "NETAPP"
        product "LUN"
        path_grouping_policy group_by_prio
        path_selector "queue-length 0"
        path_checker tur
        features "2 pg_init_retries 50"
        no_path_retry 30
        hardware_handler "1 alua"
        prio alua
        failback immediate
        rr_weight uniform
        rr_min_io_rq 1
        fast_io_fail_tmo 5
        dev_loss_tmo 60
    }
}
# END jt-pve-storage-netapp
EOF

    echo "NetApp multipath configuration added successfully."

    # Reload multipathd if running
    if systemctl is-active --quiet multipathd 2>/dev/null; then
        echo "Reloading multipathd configuration..."
        timeout 10 multipathd reconfigure 2>/dev/null || true
    else
        echo "Note: multipathd is not running. Start it with: systemctl start multipathd"
    fi

    return 0
}

case "$1" in
    configure)
        echo "=============================================="
        echo "NetApp ONTAP Storage Plugin for Proxmox VE"
        echo "=============================================="
        echo ""

        # Add NetApp multipath configuration
        add_netapp_multipath_config
        echo ""

        echo "To add storage, use:"
        echo ""
        echo "  pvesm add netappontap <storage-id> \\"
        echo "    --ontap-portal <ontap-mgmt-ip> \\"
        echo "    --ontap-svm <svm-name> \\"
        echo "    --ontap-aggregate <aggregate-name> \\"
        echo "    --ontap-username <api-user> \\"
        echo "    --ontap-password <api-password> \\"
        echo "    --content images --shared 1"
        echo ""

        # Ensure required services are enabled (with timeout to prevent hang)
        echo "Checking required services..."

        if timeout 5 systemctl is-active --quiet iscsid 2>/dev/null; then
            echo "  [OK] iSCSI initiator daemon (iscsid) is running."
        else
            echo "  [WARN] iscsid is not running. Starting..."
            timeout 10 systemctl enable iscsid 2>/dev/null || true
            timeout 10 systemctl start iscsid 2>/dev/null || echo "  [FAIL] Could not start iscsid"
        fi

        if timeout 5 systemctl is-active --quiet multipathd 2>/dev/null; then
            echo "  [OK] Multipath daemon (multipathd) is running."
        else
            echo "  [WARN] multipathd is not running. Starting..."
            timeout 10 systemctl enable multipathd 2>/dev/null || true
            timeout 10 systemctl start multipathd 2>/dev/null || echo "  [FAIL] Could not start multipathd"
        fi

        # Restart PVE services to load the new storage plugin
        echo "Restarting PVE services to load plugin..."
        if timeout 30 systemctl restart pvedaemon 2>/dev/null; then
            echo "  [OK] pvedaemon restarted."
        else
            echo "  [WARN] Could not restart pvedaemon. Please run: systemctl restart pvedaemon"
        fi

        if timeout 30 systemctl restart pveproxy 2>/dev/null; then
            echo "  [OK] pveproxy restarted."
        else
            echo "  [WARN] Could not restart pveproxy. Please run: systemctl restart pveproxy"
        fi

        echo ""
        echo "Installation complete."
        echo "See /usr/share/doc/jt-pve-storage-netapp/ for documentation."
        echo "=============================================="
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
