# debian/postinst
#!/bin/bash
set -e

# Marker comment to identify our configuration
MULTIPATH_MARKER_BEGIN="# BEGIN jt-pve-storage-netapp"
MULTIPATH_MARKER_END="# END jt-pve-storage-netapp"
MULTIPATH_CONF="/etc/multipath.conf"

# Function to add NetApp multipath configuration
add_netapp_multipath_config() {
    # Ensure /etc exists (should always exist, but be safe)
    if [ ! -d "/etc" ]; then
        echo "ERROR: /etc directory does not exist!"
        return 1
    fi

    # Check if our configuration already exists
    if [ -f "$MULTIPATH_CONF" ] && grep -q "$MULTIPATH_MARKER_BEGIN" "$MULTIPATH_CONF" 2>/dev/null; then
        echo "NetApp multipath configuration already exists in $MULTIPATH_CONF"
        return 0
    fi

    # Check if any NetApp device configuration already exists (manual setup)
    if [ -f "$MULTIPATH_CONF" ] && grep -q 'vendor.*"NETAPP"' "$MULTIPATH_CONF" 2>/dev/null; then
        echo "Existing NetApp configuration found in $MULTIPATH_CONF, skipping auto-configuration."
        echo "Note: Manual configuration will be preserved."
        echo ""
        # Check for dangerous infinite-queue settings
        local needs_update=0
        if grep -q 'no_path_retry.*queue' "$MULTIPATH_CONF" 2>/dev/null; then
            needs_update=1
        fi
        if grep -q 'queue_if_no_path' "$MULTIPATH_CONF" 2>/dev/null; then
            needs_update=1
        fi
        if grep -q 'dev_loss_tmo.*infinity' "$MULTIPATH_CONF" 2>/dev/null; then
            needs_update=1
        fi
        if [ "$needs_update" -eq 1 ]; then
            echo ""
            echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
            echo -e "\033[1;31m  WARNING: DANGEROUS MULTIPATH SETTINGS DETECTED\033[0m"
            echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
            echo ""
            echo -e "  Your \033[1m$MULTIPATH_CONF\033[0m contains settings that may cause"
            echo -e "  the \033[1;31mENTIRE PVE NODE TO HANG\033[0m if a NetApp LUN becomes unavailable."
            echo ""
            echo -e "  \033[1mRequired changes:\033[0m"
            echo ""
            echo -e "    \033[1;31mno_path_retry queue\033[0m    -->  \033[1;32mno_path_retry 30\033[0m"
            echo -e "    \033[1;31mqueue_if_no_path\033[0m       -->  \033[1;32m(remove, use no_path_retry instead)\033[0m"
            echo -e "    \033[1;31mdev_loss_tmo infinity\033[0m  -->  \033[1;32mdev_loss_tmo 60\033[0m"
            echo ""
            echo -e "  \033[1mAlso add if not present:\033[0m"
            echo -e "    \033[1;32mfast_io_fail_tmo 5\033[0m"
            echo ""
            echo -e "  \033[1mAfter editing multipath.conf:\033[0m"
            echo -e "    \033[1;32msystemctl restart multipathd\033[0m   (applies new config and flushes stale maps)"
            echo ""
            echo -e "  \033[1;33mIMPORTANT:\033[0m 'systemctl reload multipathd' only re-reads config but does"
            echo -e "  NOT clean up existing stale multipath maps. You must use 'restart' to"
            echo -e "  apply new settings to existing devices and remove orphaned entries."
            echo ""
            echo -e "  \033[1;31mWARNING:\033[0m Do NOT use 'multipath -F' (capital F) to clean stale maps."
            echo -e "  It flushes ALL unused maps, which may include working maps from other"
            echo -e "  storage systems. Instead, target specific stale WWIDs:"
            echo -e "    \033[1;32mmultipath -f <wwid>\033[0m   (lowercase f, removes one specific map)"
            echo ""
            echo -e "  For details, see: \033[1;36mhttps://github.com/jasoncheng7115/jt-pve-storage-netapp\033[0m"
            echo -e "  Refer to \033[1mdocs/CONFIGURATION.md\033[0m -> Multipath Configuration section."
            echo ""
            echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
            echo ""
        fi
        return 0
    fi

    echo "Adding NetApp multipath configuration to $MULTIPATH_CONF..."

    # Create backup if file exists
    if [ -f "$MULTIPATH_CONF" ]; then
        cp "$MULTIPATH_CONF" "${MULTIPATH_CONF}.bak.$(date +%Y%m%d%H%M%S)"
        echo "Backup created: ${MULTIPATH_CONF}.bak.$(date +%Y%m%d%H%M%S)"
    fi

    # If multipath.conf doesn't exist, create with defaults
    if [ ! -f "$MULTIPATH_CONF" ]; then
        echo "Creating new $MULTIPATH_CONF..."
        cat > "$MULTIPATH_CONF" << 'DEFAULTS_EOF'
# Multipath configuration file
# Created by jt-pve-storage-netapp package

defaults {
    user_friendly_names yes
    find_multipaths yes
}

DEFAULTS_EOF
    fi

    # Add our configuration
    cat >> "$MULTIPATH_CONF" << 'EOF'

# BEGIN jt-pve-storage-netapp
# NetApp ONTAP device configuration
# Added automatically by jt-pve-storage-netapp package
# This section will be removed when the package is purged
devices {
    device {
        vendor "NETAPP"
        product "LUN"
        path_grouping_policy group_by_prio
        path_selector "queue-length 0"
        path_checker tur
        features "2 pg_init_retries 50"
        no_path_retry 30
        hardware_handler "1 alua"
        prio alua
        failback immediate
        rr_weight uniform
        rr_min_io_rq 1
        fast_io_fail_tmo 5
        dev_loss_tmo 60
    }
}
# END jt-pve-storage-netapp
EOF

    echo "NetApp multipath configuration added successfully."

    # Reload multipathd if running
    if systemctl is-active --quiet multipathd 2>/dev/null; then
        echo "Reloading multipathd configuration..."
        timeout 10 multipathd reconfigure 2>/dev/null || true
    else
        echo "Note: multipathd is not running. Start it with: systemctl start multipathd"
    fi

    return 0
}

# Scan for stale NETAPP multipath devices (all paths failed) that may have
# been left over from previous plugin versions before automatic orphan
# cleanup was added (v0.2.2+). These devices can cause node hangs if any
# process touches them, especially with queue_if_no_path settings.
#
# We do NOT auto-clean these because we cannot reliably distinguish them
# from manually-configured NetApp storage. Instead, we list them with the
# exact commands the operator should run.
detect_stale_netapp_devices() {
    if ! command -v multipath >/dev/null 2>&1; then
        return 0
    fi

    # Get multipath -ll output, capture WWIDs of NETAPP devices with all paths failed
    local mpath_output
    mpath_output=$(timeout 15 multipath -ll 2>/dev/null) || return 0

    # Parse: find blocks for NETAPP vendor where all paths are "failed faulty"
    local stale_wwids
    stale_wwids=$(echo "$mpath_output" | awk '
        /NETAPP/ {
            in_netapp=1
            # Try to extract WWID from current or previous line
            for (i=1;i<=NF;i++) {
                if ($i ~ /^[0-9a-f]{32,}$/) wwid=$i
                if ($i ~ /^\([0-9a-f]{32,}\)$/) { wwid=$i; gsub(/[()]/,"",wwid) }
            }
            # Also check first field which is often the WWID itself
            if ($1 ~ /^[0-9a-f]{32,}$/) wwid=$1
            has_active=0
            has_failed=0
            next
        }
        /^[^ \t]/ && in_netapp {
            # New top-level entry, finalize previous
            if (has_failed && !has_active && wwid) print wwid
            in_netapp=0
            wwid=""
            has_active=0
            has_failed=0
            # Re-check current line
            if ($0 ~ /NETAPP/) {
                in_netapp=1
                for (i=1;i<=NF;i++) {
                    if ($i ~ /^[0-9a-f]{32,}$/) wwid=$i
                    if ($i ~ /^\([0-9a-f]{32,}\)$/) { wwid=$i; gsub(/[()]/,"",wwid) }
                }
                if ($1 ~ /^[0-9a-f]{32,}$/) wwid=$1
            }
            next
        }
        in_netapp && /active ready running/ { has_active=1 }
        in_netapp && /failed faulty/ { has_failed=1 }
        END {
            if (has_failed && !has_active && wwid) print wwid
        }
    ')

    if [ -z "$stale_wwids" ]; then
        return 0
    fi

    local count
    count=$(echo "$stale_wwids" | wc -l)

    echo ""
    echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
    echo -e "\033[1;31m  STALE NETAPP MULTIPATH DEVICES DETECTED ($count)\033[0m"
    echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
    echo ""
    echo "  These multipath devices have all paths in failed state. They may be"
    echo "  leftover from a previous plugin version that did not fully clean up."
    echo ""
    echo "  These devices are DANGEROUS:"
    echo "  - Any process touching them (e.g., 'vgs', 'lvs', 'lsblk') may hang"
    echo "  - With 'queue_if_no_path' multipath setting, the entire node can hang"
    echo "  - Migration tasks may fail with cascading errors"
    echo ""
    echo -e "  \033[1;31mPlease verify each WWID below is NOT in use by other storage,\033[0m"
    echo -e "  \033[1;31mthen clean them manually with:\033[0m"
    echo ""
    for wwid in $stale_wwids; do
        echo -e "    \033[1;32mmultipath -f $wwid\033[0m"
    done
    echo ""
    echo "  After flushing, also remove stale SCSI device entries:"
    echo -e "    \033[1;32mfor sd in \$(lsscsi | grep NETAPP | awk '{print \$NF}' | xargs -I{} basename {}); do\033[0m"
    echo -e "    \033[1;32m  echo 1 > /sys/block/\$sd/device/delete\033[0m"
    echo -e "    \033[1;32mdone\033[0m"
    echo ""
    echo -e "  \033[1;33mWARNING:\033[0m If 'multipath -f' hangs (queue_if_no_path), use:"
    echo -e "    \033[1;32mmultipathd disablequeueing map <wwid>\033[0m"
    echo -e "    \033[1;32mdmsetup message <wwid> 0 fail_if_no_path\033[0m"
    echo -e "    \033[1;32mdmsetup remove --force --retry <wwid>\033[0m"
    echo ""
    echo -e "  \033[1;31mDO NOT use 'multipath -F' (capital F)\033[0m -- it will flush ALL"
    echo "  unused maps including manually-configured storage on idle nodes."
    echo ""
    echo "  See: https://github.com/jasoncheng7115/jt-pve-storage-netapp"
    echo "  Doc: docs/TROUBLESHOOTING.md -> Stale Multipath Devices"
    echo ""
    echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
    echo ""
}

case "$1" in
    configure)
        echo "=============================================="
        echo "NetApp ONTAP Storage Plugin for Proxmox VE"
        echo "=============================================="
        echo ""

        # Add NetApp multipath configuration
        add_netapp_multipath_config
        echo ""

        # Detect stale multipath devices left over from previous versions
        detect_stale_netapp_devices

        # Check if lvm.conf global_filter excludes plugin multipath devices.
        # Without this filter, the host's LVM scanner reads VG metadata from
        # inside VM disks on plugin-managed LUNs and auto-activates them as
        # host-level VGs. This creates holders on the multipath device that
        # block free_image() from deleting volumes. Common on PVE nodes
        # upgraded from 7->8->9 that inherited old lvm.conf without filter.
        LVM_CONF="/etc/lvm/lvm.conf"
        if [ -f "$LVM_CONF" ]; then
            if ! grep -qE '^\s*global_filter\s*=' "$LVM_CONF" 2>/dev/null; then
                echo ""
                echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
                echo -e "\033[1;33m  WARNING: /etc/lvm/lvm.conf has no global_filter\033[0m"
                echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
                echo ""
                echo "  Without a global_filter, the host's LVM scanner will read VG metadata"
                echo "  from inside VM disks on plugin-managed multipath LUNs and auto-activate"
                echo "  them as host-level VGs. This causes:"
                echo "    - 'Cannot delete volume: device is still in use' errors"
                echo "    - move-disk unable to clean up source disk after migration"
                echo "    - Duplicate VG name warnings from vgs/pvs"
                echo ""
                echo "  This is common on PVE nodes upgraded from 7 -> 8 -> 9."
                echo "  Fresh PVE 9 installs have a correct global_filter by default."
                echo ""
                echo -e "  \033[1;32mRecommended fix:\033[0m Add to $LVM_CONF devices { } section:"
                echo ""
                echo -e "    \033[1;32mglobal_filter = [ \"r|/dev/mapper/360.*|\", \"r|/dev/dm-|\", \"a|.*|\" ]\033[0m"
                echo ""
                echo "  This excludes plugin-managed multipath devices (WWID starting with 360)"
                echo "  and raw dm-* devices from LVM scanning, while allowing everything else."
                echo "  Adjust the regex if you have non-plugin multipath devices that LVM should see."
                echo ""
                echo "  After editing, run: vgscan --mknodes"
                echo ""
                echo -e "\033[1;33m!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!\033[0m"
                echo ""
            fi
        fi

        echo "To add storage, use:"
        echo ""
        echo "  pvesm add netappontap <storage-id> \\"
        echo "    --ontap-portal <ontap-mgmt-ip> \\"
        echo "    --ontap-svm <svm-name> \\"
        echo "    --ontap-aggregate <aggregate-name> \\"
        echo "    --ontap-username <api-user> \\"
        echo "    --ontap-password <api-password> \\"
        echo "    --content images --shared 1"
        echo ""

        # Ensure required services are enabled (with timeout to prevent hang)
        echo "Checking required services..."

        if timeout 5 systemctl is-active --quiet iscsid 2>/dev/null; then
            echo "  [OK] iSCSI initiator daemon (iscsid) is running."
        else
            echo "  [WARN] iscsid is not running. Starting..."
            timeout 10 systemctl enable iscsid 2>/dev/null || true
            timeout 10 systemctl start iscsid 2>/dev/null || echo "  [FAIL] Could not start iscsid"
        fi

        if timeout 5 systemctl is-active --quiet multipathd 2>/dev/null; then
            echo "  [OK] Multipath daemon (multipathd) is running."
        else
            echo "  [WARN] multipathd is not running. Starting..."
            timeout 10 systemctl enable multipathd 2>/dev/null || true
            timeout 10 systemctl start multipathd 2>/dev/null || echo "  [FAIL] Could not start multipathd"
        fi

        # Reload PVE services to load the new storage plugin code.
        #
        # Use 'reload' (SIGHUP) instead of 'restart' when the service is
        # already running. PVE::Daemon handles SIGHUP by re-exec'ing itself
        # with the same PID, picking up new Perl modules from disk without
        # going through a stop phase. This avoids the bootstrapping problem
        # where the OLD code has already created D-state children (e.g. from
        # writing to a non-iSCSI SCSI host's scan file on HPE smartpqi),
        # and 'systemctl restart' hangs waiting for those unkillable children
        # during the stop phase.
        #
        # CRITICAL: pvestatd must also be reloaded. It polls status() every
        # 10 seconds using whatever plugin code it loaded at startup. If we
        # only reload pvedaemon+pveproxy but leave pvestatd running old code,
        # pvestatd continues creating D-state children from old rescan logic.
        # This was the root cause of a production watchdog-triggered reboot
        # (v0.2.5 incident: pvestatd with v0.2.4 code kept hitting host1).
        echo "Reloading PVE services to load plugin..."
        for svc in pvedaemon pvestatd pveproxy; do
            if timeout 5 systemctl is-active --quiet "$svc" 2>/dev/null; then
                # Service is running -> reload (SIGHUP, re-exec with new code)
                if timeout 15 systemctl reload "$svc" 2>/dev/null; then
                    echo "  [OK] $svc reloaded."
                else
                    echo "  [WARN] Could not reload $svc. Please run: systemctl reload $svc"
                fi
            else
                # Service is not running -> start fresh
                if timeout 30 systemctl start "$svc" 2>/dev/null; then
                    echo "  [OK] $svc started."
                else
                    echo "  [WARN] Could not start $svc. Please run: systemctl start $svc"
                fi
            fi
        done

        echo ""
        echo "Installation complete."
        echo "See /usr/share/doc/jt-pve-storage-netapp/ for documentation."
        echo "=============================================="
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
